# rc.d/dhcpleasemon
#!/bin/ksh

daemon="/usr/local/sbin/dhcpleasemon"

. /etc/rc.d/rc.subr

rc_cmd $1
